# kde-cachy/tether-daemon.service
[Unit]
Description=Tether Network Mount Daemon (KDE/CachyOS)
Documentation=https://github.com/InsaneAI1973/Tether
After=graphical-session.target plasma-ksmserver.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/tether-kde/daemon.py
Restart=on-failure
RestartSec=5
Environment=PYTHONPATH=/opt/tether-kde
Environment=DBUS_SESSION_BUS_ADDRESS=%I
Environment=QT_LOGGING_RULES=qt.qpa.wayland=false

[Install]
WantedBy=graphical-session.target
